# config/supervisor.conf
[program:orchestrator]
command=/usr/local/bin/python3 main.py
directory=/opt/app
user=orchestrator
autostart=true
autorestart=true
stderr_logfile=/var/log/orchestrator/err.log
stdout_logfile=/var/log/orchestrator/out.log

[program:webhook]
command=/usr/local/bin/python3 services/webhook_server.py
directory=/opt/app
user=orchestrator
autostart=true
autorestart=true